# systemd/sec-journal.service
[Unit]
Description=sec daemon for journal
After=basic.target

[Service]
ExecStart=/usr/lib/arch-setup/sec.sh

[Install]
WantedBy=multi-user.target
